# ansible/roles/qortex/templates/qortex-mcp.service.j2
# {{ ansible_managed }}
# qortex MCP server — Streamable HTTP transport for gateway connections
# Runs alongside qortex.service (REST API) on a separate port.

[Unit]
Description=Qortex MCP Server (Streamable HTTP)
After=network-online.target qortex.service
Wants=network-online.target

[Service]
Type=simple
User={{ ansible_user }}
Group={{ ansible_user }}
EnvironmentFile=/etc/openclaw/qortex.env
{% if qortex_otel_enabled | default(false) | bool %}
EnvironmentFile=/etc/openclaw/qortex-otel.env
{% endif %}
# FastMCP reads FASTMCP_HOST / FASTMCP_PORT from env
Environment=FASTMCP_HOST={{ qortex_mcp_host }}
Environment=FASTMCP_PORT={{ qortex_mcp_port }}
ExecStart={{ user_home }}/.local/bin/qortex mcp-serve --transport streamable-http
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

# Hardening — same as qortex.service
NoNewPrivileges=true
ProtectSystem=strict
BindReadOnlyPaths={{ user_home }}/.local
BindPaths={{ user_home }}/.cache
BindPaths={{ user_home }}/.qortex

[Install]
WantedBy=multi-user.target
