# ansible/roles/qortex/templates/qortex.service.j2
# {{ ansible_managed }}
# qortex HTTP service — persistent vector search + knowledge graph backend

[Unit]
Description=Qortex Knowledge Graph Service
After=network-online.target{% if pgvector_enabled | default(false) | bool %} docker.service{% endif %}

Wants=network-online.target

[Service]
Type=simple
User={{ ansible_user }}
Group={{ ansible_user }}
EnvironmentFile=/etc/openclaw/qortex.env
{% if qortex_otel_enabled | default(false) | bool %}
EnvironmentFile=/etc/openclaw/qortex-otel.env
{% endif %}
ExecStart={{ user_home }}/.local/bin/qortex serve --host {{ qortex_serve_host }} --port {{ qortex_serve_port }}
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

# Hardening — ProtectHome=read-only conflicts with ReadWritePaths on some
# systemd versions (the model cache under ~/.cache needs write access).
# BindPaths expose only the directories the service actually needs.
NoNewPrivileges=true
ProtectSystem=strict
BindReadOnlyPaths={{ user_home }}/.local
BindPaths={{ user_home }}/.cache
BindPaths={{ user_home }}/.qortex

[Install]
WantedBy=multi-user.target
